# madsonic.conf
[program:start]
user = nobody
command = /var/madsonic/start.sh